# deploy/autodl/supervisord.conf
[supervisord]
nodaemon=false
logfile=%(ENV_CRA_DATA_ROOT)s/logs/supervisord.log
pidfile=%(ENV_CRA_DATA_ROOT)s/run/supervisord.pid
childlogdir=%(ENV_CRA_DATA_ROOT)s/logs

[program:inference]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh inference
autostart=true
autorestart=true
priority=10
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/inference.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/inference.err.log

[program:api]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh api
autostart=true
autorestart=true
priority=20
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/api.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/api.err.log

[program:dispatcher]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh dispatcher
autostart=true
autorestart=true
priority=30
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/dispatcher.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/dispatcher.err.log

[program:worker-analysis]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh worker-analysis
autostart=true
autorestart=true
priority=40
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-analysis.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-analysis.err.log

[program:worker-research]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh worker-research
autostart=true
autorestart=true
priority=40
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-research.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-research.err.log

[program:worker-evaluation]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh worker-evaluation
autostart=true
autorestart=true
priority=40
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-evaluation.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-evaluation.err.log

[program:worker-maintenance]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh worker-maintenance
autostart=true
autorestart=true
priority=40
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-maintenance.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/worker-maintenance.err.log

[program:beat]
command=%(ENV_CRA_REPO_ROOT)s/deploy/autodl/run-service.sh beat
autostart=true
autorestart=true
priority=50
stdout_logfile=%(ENV_CRA_DATA_ROOT)s/logs/beat.log
stderr_logfile=%(ENV_CRA_DATA_ROOT)s/logs/beat.err.log
